# scripts/com.ballast.agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.ballast.agent</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/mac/Vibecoding/Ballast/scripts/run_tick.sh</string>
  </array>

  <!-- Run once an hour (the agent's main clock). -->
  <key>StartInterval</key>
  <integer>3600</integer>

  <!-- Run one tick immediately when loaded. -->
  <key>RunAtLoad</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/mac/Vibecoding/Ballast/ballast-agent.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/mac/Vibecoding/Ballast/ballast-agent.log</string>

  <key>WorkingDirectory</key>
  <string>/Users/mac/Vibecoding/Ballast</string>

  <!-- Don't relaunch on its own; StartInterval handles cadence. -->
  <key>KeepAlive</key>
  <false/>
</dict>
</plist>
